# backend/scripts/agc-backend.service
# AGC backend — systemd unit
# This mirrors the unit ACTUALLY RUNNING in production as of AGC-045.2
# (confirmed 2026-07-03: venv at repo root, not backend/venv; Restart=always).
# It exists here as a versioned reference/backup only — do not install this
# over the live unit without an explicit decision to do so.
#
# AGC-071: binds 127.0.0.1, matching docs/deploy.md and nginx/agc.conf
# (which already proxy_pass to 127.0.0.1:8000). The previous 0.0.0.0 bind
# relied entirely on ufw denying external access to 8000 as its only
# defense — binding loopback-only removes that single point of failure.

[Unit]
Description=AGC AI Highlight Creator - FastAPI backend
After=network.target

[Service]
Type=simple
User=agc
WorkingDirectory=/home/agc/agc-ai-highlight-creator/backend
Environment=PATH=/home/agc/agc-ai-highlight-creator/venv/bin
ExecStart=/home/agc/agc-ai-highlight-creator/venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
